# services/halloween-video.service
[Unit]
Description=Halloween Unified Video Player Service
After=network.target graphical-session.target
Wants=graphical-session.target

[Service]
Type=simple
User=volvo
Group=volvo
WorkingDirectory=/home/volvo/Desktop/Development/first_robotics_falcon_halloween/lib/rpi
ExecStart=/usr/bin/python3 /home/volvo/Desktop/Development/first_robotics_falcon_halloween/lib/rpi/rpi_video_screen.py single
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Timeout settings for graceful shutdown (reduced from 15s since VLC shuts down faster than OpenCV)
TimeoutStopSec=5
KillMode=process

# Environment variables for display (required for VLC video output)
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/volvo/.Xauthority

# Ensure the service has access to GPIO and audio
SupplementaryGroups=gpio audio video

[Install]
WantedBy=graphical.target